# build/shared/supervisor/octane-worker-dev.conf
[program:octane]
process_name=%(program_name)s
command=php /app/artisan octane:start --workers=4 --watch --no-interaction --port=8080
autostart=true
autorestart=true
stdout_logfile=/app/storage/logs/supervisor/octane.log
stderr_logfile=/app/storage/logs/supervisor/octane-error.log
stopwaitsecs=3600
