# src/shared/supervisor/octane-worker-dev.conf
[program:laravel-octane]
process_name=%(program_name)s
command=php /app/artisan octane:start --workers=4 --watch --no-interaction --port=8080
autostart=true
autorestart=true
stopwaitsecs=3600
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
